# systemd/zeusnet-serial.service
[Unit]
Description=ZeusNet Serial Command Bus
After=network.target

[Service]
User=zeus
WorkingDirectory=/opt/zeusnet
EnvironmentFile=/opt/zeusnet/.env
ExecStart=/opt/zeusnet/.venv/bin/python backend/c2/command_bus.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
